# overlay.d/05core/usr/lib/dracut/modules.d/35coreos-ignition/coreos-gpt-setup.service
[Unit]
Description=Generate New UUID For Boot Disk GPT
ConditionPathExists=/etc/initrd-release
Before=systemd-fsck-root.service
Before=ignition-diskful.target

# Have some ordering with this service since it also wants to access the bootfs.
After=coreos-copy-firstboot-network.service

Requires=dev-disk-by\x2dlabel-boot.device
After=dev-disk-by\x2dlabel-boot.device

# And since the boot device may be on multipath; optionally wait for it to
# appear via the dynamic target.
After=coreos-unique-boot.service
Before=ignition-ostree-uuid-boot.service

# Run before services that use device nodes, preventing them from racing
# with udev activity generated by sgdisk
Before=ignition-kargs.service

OnFailure=emergency.target
OnFailureJobMode=isolate

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/coreos-gpt-setup /dev/disk/by-label/boot
