# overlay.d/05core/usr/lib/systemd/system/emergency.service.d/coreos-sulogin-force.conf
# https://github.com/coreos/coreos-installer/commit/15a79263d0bd5d72056a6080f6687dc10cba2dda
# https://github.com/systemd/systemd/pull/10397
# We want things like `systemd.unit=emergency.target` and `single` on the
# kernel command line to just work even with our locked root account.
# This file is used as an override for both emergency.target and rescue.target.
[Service]
Environment=SYSTEMD_SULOGIN_FORCE=1
