# overlay.d/12kdump/usr/lib/systemd/system/kdump.service.d/remount-boot.conf
# https://bugzilla.redhat.com/show_bug.cgi?id=1918493
# `/boot` is read-only, but `kdump.service` wants to 
# places its generated initramfs alongside the default
# initramfs under `/boot/ostree`.
# Until `kdump` gains the ability to place its initramfs
# elsewhere, temporarily remount `/boot` read-write before
# the `kdump` initramfs is generated.
[Service]
ExecStartPre=/usr/bin/mount -o remount,rw /boot
