# overlay.d/05core/usr/lib/systemd/system/multipathd.socket.d/50-start-conditions.conf
# Temporary workaround for https://bugzilla.redhat.com/show_bug.cgi?id=2008098
[Unit]
ConditionKernelCommandLine=!multipath=off
ConditionKernelCommandLine=!nompath
ConditionPathExists=/etc/multipath.conf
ConditionVirtualization=!container
