# overlay.d/05core/usr/lib/systemd/system/ignition-delete-config.service.d/10-flag-file.conf
# Create a flag file to notify coreos-ignition-delete-config.service that
# we've run, and put it in /run because /var isn't mounted yet.
# coreos-ignition-delete-config.service will then avoid trying to delete
# the config again, and will create a persistent stamp file in /var/lib.

[Service]
ExecStart=/bin/touch /run/coreos-ignition-delete-config.stamp
