# overlay.d/05core/usr/lib/systemd/system/systemd-backlight@.service.d/45-after-ostree-remount.conf
# Temporary fix for https://github.com/coreos/fedora-coreos-tracker/issues/975
# until https://github.com/ostreedev/ostree/issues/2115 is fixed.
[Unit]
After=ostree-remount.service
