# overlay.d/15fcos/usr/lib/systemd/system/coreos-reset-stub-resolv-selinux-context.path
[Unit]
Description=Monitor /run/systemd/resolve/stub-resolv.conf to fixup SELinux context.
Documentation=https://github.com/systemd/systemd/pull/17976
Before=systemd-resolved.service

[Path]
PathModified=/run/systemd/resolve/stub-resolv.conf
Unit=coreos-reset-stub-resolv-selinux-context.service

[Install]
WantedBy=multi-user.target
