# overlay.d/05core/usr/lib/systemd/system/boot.mount
[Unit]
Description=Boot partition
Documentation=https://github.com/coreos/fedora-coreos-config

[Mount]
What=/dev/disk/by-label/boot
Where=/boot
