# overlay.d/05core/usr/lib/systemd/system/coreos-growpart.service
# See also coreos-growpart script.

[Unit]
ConditionPathExists=!/var/lib/coreos-growpart.stamp
Before=sshd.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/coreos-growpart /
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
