# overlay.d/05core/usr/lib/dracut/modules.d/50coreos-kernel/coreos-check-kernel.service
[Unit]
Description=Check That Initrd Matches Kernel
DefaultDependencies=false
Before=sysinit.target systemd-modules-load.service
ConditionPathIsDirectory=!/usr/lib/modules/%v

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/false
