# overlay.d/05core/usr/lib/dracut/modules.d/40ignition-ostree/ignition-ostree-growfs.service
[Unit]
Description=Ignition OSTree: Grow root filesystem
DefaultDependencies=false
ConditionKernelCommandLine=ostree
ConditionPathExists=!/run/ostree-live
Before=initrd-root-fs.target
After=ignition-ostree-mount-firstboot-sysroot.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/coreos-growpart /sysroot
RemainAfterExit=yes
