# overlay.d/05core/usr/lib/systemd/system/systemd-firstboot.service.d/fcos-disable.conf
# See the comment in 40-coreos-systemd.preset; we're
# keeping this even stronger disable override for now,
# but it may not really be necessary.
[Unit]
ConditionPathExists=/run/nosuchfile
