# overlay.d/05core/usr/lib/systemd/system/systemd-firstboot.service.d/fcos-disable.conf
# Disable systemd-firstboot because it conflicts with Ignition.
# In most cases this is handled via the remove-from-packages
# bits in the manifest (ignition-and-ostree.yaml), but
# we want to support overlaying builds of systemd from git.
[Unit]
ConditionPathExists=/run/nosuchfile
